# regression/setup.cron
17 11 * * *  /home/ec2-user/myaws/DataScience/regression/run.sh
