# regression/setup.cron
PATH = /home/ec2-user/anaconda3/bin/
SHELL = /bin/bash

17 1 * * *  /home/ec2-user/myaws/DataScience/regression/run.sh
17 2 * * *  python /home/ec2-user/db.py
17 3  * * *  python /home/ec2-user/myaws/GoogleAPI/gssValue.py > /home/ec2-user/fundamentals.txt
