# regression/setup.cron
17 6 * * *  /home/ec2-user/myaws/DataScience/regression/run.sh
17 7 * * *  python /home/ec2-user/db.py
17 8 * * *  python /home/ec2-user/myaws/GoogleAPI/gssValue.py > /home/ec2-user/fundamentals.txt
